# deploy/systemd/chesspuzzle.service
[Unit]
Description=ChessPuzzle Docker Compose production stack
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/deploy/chesspuzzle
# Load .env file if it exists (contains DOMAIN, secrets, etc.)
# The "-" prefix makes this optional - systemd won't fail if the file is missing
EnvironmentFile=-/home/deploy/chesspuzzle/.env

# Start the production docker-compose stack
ExecStart=/usr/bin/docker compose -f /home/deploy/chesspuzzle/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /home/deploy/chesspuzzle/docker-compose.prod.yml down

[Install]
WantedBy=multi-user.target

# Notes:
# - Copy this file to /etc/systemd/system/chesspuzzle.service and adjust paths as needed.
# - Ensure the deploy user (or whichever user runs docker) is in the docker group.
# - The unit references /home/deploy/chesspuzzle; change to the actual repository path on the host.